# deploy/ptdtt.service
[Unit]
Description=PTDTT Manager Web App
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ptdtt-manager
ExecStart=/var/www/ptdtt-manager/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 2 \
    --timeout 120 \
    --access-logfile /var/log/ptdtt/access.log \
    --error-logfile /var/log/ptdtt/error.log \
    server_flask:app
Restart=always
RestartSec=5

# Environment
Environment="PORT=5000"
# Environment="EMR_USER=VKAN"
# Environment="EMR_PASS=anmd3010"

[Install]
WantedBy=multi-user.target
